# home/autostart/chromium.desktop
[Desktop Entry]
Type=Application
Name=Chromium
Exec=/usr/bin/chromium-browser --start-fullscreen
