# home/autostart/lower-resolution.desktop
[Desktop Entry]
Type=Application
Name=Reduce resolution
Exec=xrandr --output HDMI-1 --mode 2560x1440
